# startup/geth-supervisord.conf
[supervisord]
logfile=gexp-supervisord.log

[program:eth]
stdout_logfile=gexp.log
stderr_logfile=gexp.err
command=gexp
